# data/mirrorman-refresh.timer
[Unit]
Description=Mirrorman Auto Refresh Mirrors Timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
